# CH340-fan-control.service
[Unit]
Description=CH340-fan-control

[Service]
ExecStart=/usr/local/bin/CH340-fan-control.py

[Install]
WantedBy=multi-user.target
